# server/systemd/outliner.service
[Unit]
Description=Outliner Server
After=network.target

[Service]
Type=simple
User=outliner
Group=outliner
WorkingDirectory=/srv/outliner/server
EnvironmentFile=/etc/outliner/outliner.env
ExecStart=/usr/bin/node /srv/outliner/server/dist/index.js
Restart=on-failure
RestartSec=5s
NoNewPrivileges=yes
PrivateTmp=true
PrivateDevices=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/srv/outliner

[Install]
WantedBy=multi-user.target
